# install/motd.init
#!/sbin/openrc-run
start (){
  cat /etc/motd
}
